# bootstrap/crontab
@reboot while [ ! -e /vagrant/bootstrap/start.sh ]; do sleep 1; done; sh /vagrant/bootstrap/start.sh
